# systemd/sxhkd.service
[Unit]
Description=Keybinding Daemon

[Service]
Type=exec
ExecStart=/bin/sxhkd

[Install]
WantedBy=window-manager.target

