# systemd/window-manager.target
[Unit]
Description=Full UI is present.
AllowIsolate=no
