# systemd/wm-background.service
[Unit]
Description=Displayes a background image

[Service]
Type=oneshot
ExecStart=/home/mkr/.fehbg

[Install]
WantedBy=window-manager.target
